# service/speedtest.timer
[Unit]
Description=SpeedTest

[Timer]
OnCalendar=*-*-* *:00:00
Unit=speedtest.service

[Install]
WantedBy=multi-user.target
